# devopsfetch.service
[Unit]
Description=Devopsfetch system monitoring
After=network.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/devopsfetch --monitor
